# filesystem/usr/share/gvfs/mounts/obexftp.mount
[Mount]
Type=obex
Exec=/usr/lib/gvfs/gvfsd-obexftp
AutoMount=false
