# filesystem/var/lib/dpkg/info/alsa-utils.conffiles
/etc/alsa/modprobe-post-install.d/alsa-utils
/etc/init.d/alsa-utils
