# filesystem/var/lib/dpkg/info/aptitude.list
/.
/usr
/usr/bin
/usr/bin/aptitude-run-state-bundle
/usr/bin/aptitude
/usr/bin/aptitude-create-state-bundle
/usr/share
/usr/share/man
/usr/share/man/pl
/usr/share/man/pl/man8
/usr/share/man/pl/man8/aptitude.8.gz
/usr/share/man/gl
/usr/share/man/gl/man8
/usr/share/man/gl/man8/aptitude.8.gz
/usr/share/man/fi
/usr/share/man/fi/man8
/usr/share/man/fi/man8/aptitude.8.gz
/usr/share/man/ja
/usr/share/man/ja/man8
/usr/share/man/ja/man8/aptitude.8.gz
/usr/share/man/man8
/usr/share/man/man8/aptitude.8.gz
/usr/share/man/man1
/usr/share/man/man1/aptitude-create-state-bundle.1.gz
/usr/share/man/man1/aptitude-run-state-bundle.1.gz
/usr/share/man/fr
/usr/share/man/fr/man8
/usr/share/man/fr/man8/aptitude.8.gz
/usr/share/man/de
/usr/share/man/de/man8
/usr/share/man/de/man8/aptitude.8.gz
/usr/share/man/cs
/usr/share/man/cs/man8
/usr/share/man/cs/man8/aptitude.8.gz
/usr/share/man/it
/usr/share/man/it/man8
/usr/share/man/it/man8/aptitude.8.gz
/usr/share/aptitude
/usr/share/aptitude/aptitude-defaults.ru
/usr/share/aptitude/help-nb.txt
/usr/share/aptitude/aptitude-defaults.fr
/usr/share/aptitude/aptitude-defaults.it
/usr/share/aptitude/README.fr
/usr/share/aptitude/help-zh_TW.txt
/usr/share/aptitude/FAQ
/usr/share/aptitude/README
/usr/share/aptitude/aptitude-defaults.fi
/usr/share/aptitude/aptitude-defaults.bs
/usr/share/aptitude/mine-help-fi.txt
/usr/share/aptitude/help-fi.txt
/usr/share/aptitude/aptitude-defaults.tl
/usr/share/aptitude/aptitude-defaults.ar
/usr/share/aptitude/aptitude-defaults.pl
/usr/share/aptitude/help-ja.txt
/usr/share/aptitude/aptitude-defaults.gl
/usr/share/aptitude/section-descriptions
/usr/share/aptitude/help-fr.txt
/usr/share/aptitude/help-sk.txt
/usr/share/aptitude/mine-help-sv.txt
/usr/share/aptitude/help-de.txt
/usr/share/aptitude/aptitude-defaults
/usr/share/aptitude/aptitude-defaults.pt
/usr/share/aptitude/function_pkgs
/usr/share/aptitude/mine-help-es.txt
/usr/share/aptitude/aptitude-defaults.hu
/usr/share/aptitude/aptitude-defaults.lt
/usr/share/aptitude/help-cs.txt
/usr/share/aptitude/aptitude-defaults.nb
/usr/share/aptitude/mine-help-gl.txt
/usr/share/aptitude/README.fi
/usr/share/aptitude/aptitude-defaults.zh_TW
/usr/share/aptitude/mine-help-de.txt
/usr/share/aptitude/help-eu.txt
/usr/share/aptitude/help-gl.txt
/usr/share/aptitude/aptitude-defaults.ja
/usr/share/aptitude/aptitude-defaults.de
/usr/share/aptitude/aptitude-defaults.es
/usr/share/aptitude/help-sv.txt
/usr/share/aptitude/help-ru.txt
/usr/share/aptitude/aptitude-defaults.pt_BR
/usr/share/aptitude/mine-help-pt_BR.txt
/usr/share/aptitude/aptitude-defaults.ku
/usr/share/aptitude/aptitude-defaults.mr
/usr/share/aptitude/aptitude-defaults.el
/usr/share/aptitude/help-it.txt
/usr/share/aptitude/aptitude-defaults.km
/usr/share/aptitude/help-pl.txt
/usr/share/aptitude/aptitude-defaults.ro
/usr/share/aptitude/aptitude-defaults.nl
/usr/share/aptitude/help-tr.txt
/usr/share/aptitude/aptitude-defaults.cs
/usr/share/aptitude/aptitude-defaults.sk
/usr/share/aptitude/mine-help-cs.txt
/usr/share/aptitude/mine-help-ja.txt
/usr/share/aptitude/help-uk.txt
/usr/share/aptitude/help-es.txt
/usr/share/aptitude/aptitude-defaults.eu
/usr/share/aptitude/help.txt
/usr/share/aptitude/mine-help.txt
/usr/share/aptitude/function_groups
/usr/share/aptitude/mine-help-fr.txt
/usr/share/aptitude/aptitude-defaults.ca
/usr/share/aptitude/mine-help-sk.txt
/usr/share/aptitude/aptitude-defaults.sv
/usr/share/aptitude/aptitude-defaults.da
/usr/share/aptitude/help-zh_CN.txt
/usr/share/aptitude/aptitude-defaults.zh_CN
/usr/share/aptitude/help-pt_BR.txt
/usr/share/aptitude/aptitude-defaults.nn
/usr/share/aptitude/aptitude-defaults.tr
/usr/share/aptitude/aptitude-defaults.vi
/usr/share/aptitude/mine-help-it.txt
/usr/share/aptitude/NEWS
/usr/share/doc
/usr/share/doc/aptitude
/usr/share/doc/aptitude/README.Debian
/usr/share/doc/aptitude/examples
/usr/share/doc/aptitude/examples/pkg_hier_dump.cc.gz
/usr/share/doc/aptitude/examples/pkg_hier.h.gz
/usr/share/doc/aptitude/examples/pkg_hier.cc.gz
/usr/share/doc/aptitude/changelog.Debian.gz
/usr/share/doc/aptitude/README.hier
/usr/share/doc/aptitude/copyright
/usr/share/doc/aptitude/TODO
/usr/share/menu
/usr/share/menu/aptitude
/usr/share/bug
/usr/share/bug/aptitude
/var
/var/lib
/var/lib/aptitude
/etc
/etc/logrotate.d
/etc/logrotate.d/aptitude
/etc/cron.daily
/etc/cron.daily/aptitude
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/05aptitude
/usr/share/aptitude/COPYING
/usr/share/doc/aptitude/README
/usr/share/doc/aptitude/NEWS
