# filesystem/var/lib/dpkg/info/gdm.list
/.
/var
/var/lib
/var/lib/gdm
/var/log
/var/log/gdm
/usr
/usr/bin
/usr/bin/gdmthemetester
/usr/bin/gdmflexiserver
/usr/bin/gdm-dmx-reconnect-proxy
/usr/bin/gdmphotosetup
/usr/bin/gdmXnestchooser
/usr/bin/gdmdynamic
/usr/sbin
/usr/sbin/gdm
/usr/sbin/gdmsetup
/usr/lib
/usr/lib/gdm
/usr/lib/gdm/gdmchooser
/usr/lib/gdm/gdmprefetch
/usr/lib/gdm/gdmaskpass
/usr/lib/gdm/gdmgreeter
/usr/lib/gdm/gdm-ssh-session
/usr/lib/gdm/gdmopen
/usr/lib/gdm/gdmtranslate
/usr/lib/gdm/gdmlogin
/usr/lib/gdmplay
/usr/lib/gtk-2.0
/usr/lib/gtk-2.0/modules
/usr/lib/gtk-2.0/modules/libkeymouselistener.so
/usr/lib/gtk-2.0/modules/libdwellmouselistener.so
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/gdmflexiserver.1.gz
/usr/share/man/man1/gdm.1.gz
/usr/share/man/man8
/usr/share/man/man8/gdmchooser.8.gz
/usr/share/man/man8/gdmlogin.8.gz
/usr/share/man/man8/gdm.8.gz
/usr/share/pixmaps
/usr/share/pixmaps/gdm.xpm
/usr/share/pixmaps/nobody.png
/usr/share/pixmaps/nohost.png
/usr/share/pixmaps/gdmDebianLogo.xpm
/usr/share/pixmaps/login-photo.xpm
/usr/share/pixmaps/gdm-foot-logo.png
/usr/share/pixmaps/faces
/usr/share/pixmaps/faces/sunflower.jpg
/usr/share/pixmaps/faces/sunset.jpg
/usr/share/pixmaps/faces/coffee.jpg
/usr/share/pixmaps/faces/flower.jpg
/usr/share/pixmaps/faces/launch.jpg
/usr/share/pixmaps/faces/fish.jpg
/usr/share/pixmaps/faces/sky.jpg
/usr/share/pixmaps/faces/flake.jpg
/usr/share/pixmaps/faces/dice.jpg
/usr/share/pixmaps/faces/leaf.jpg
/usr/share/pixmaps/faces/butterfly.png
/usr/share/pixmaps/faces/penguin.jpg
/usr/share/pixmaps/faces/grapes.jpg
/usr/share/pixmaps/faces/baseball.png
/usr/share/pixmaps/faces/energy-arc.jpg
/usr/share/pixmaps/faces/soccerball.png
/usr/share/pixmaps/faces/lightning.jpg
/usr/share/pixmaps/faces/chess.jpg
/usr/share/pixmaps/faces/puppy.jpg
/usr/share/pixmaps/faces/guitar.jpg
/usr/share/pixmaps/faces/tennis-ball.png
/usr/share/pixmaps/faces/yellow-rose.jpg
/usr/share/pixmaps/faces/cat-eye.jpg
/usr/share/pixmaps/faces/astronaut.jpg
/usr/share/gdm
/usr/share/gdm/gdmchooser.glade
/usr/share/gdm/themes
/usr/share/gdm/themes/happygnome-list
/usr/share/gdm/themes/happygnome-list/happygnome-list.xml
/usr/share/gdm/themes/happygnome-list/background.svg
/usr/share/gdm/themes/happygnome-list/GdmGreeterTheme.desktop
/usr/share/gdm/themes/happygnome-list/gnome-logo.svg
/usr/share/gdm/themes/happygnome-list/screenshot.png
/usr/share/gdm/themes/circles
/usr/share/gdm/themes/circles/background.svg
/usr/share/gdm/themes/circles/GdmGreeterTheme.desktop
/usr/share/gdm/themes/circles/circles.xml
/usr/share/gdm/themes/circles/screenshot.png
/usr/share/gdm/themes/circles/flower.png
/usr/share/gdm/themes/happygnome
/usr/share/gdm/themes/happygnome/happygnome.xml
/usr/share/gdm/themes/happygnome/background.svg
/usr/share/gdm/themes/happygnome/GdmGreeterTheme.desktop
/usr/share/gdm/themes/happygnome/gnome-logo.svg
/usr/share/gdm/themes/happygnome/screenshot.png
/usr/share/gdm/gdmXnestWrapper
/usr/share/gdm/gdmsetup.glade
/usr/share/gdm/gdmphotosetup.glade
/usr/share/gdm/BuiltInSessions
/usr/share/gdm/BuiltInSessions/default.desktop
/usr/share/omf
/usr/share/omf/gdm
/usr/share/omf/gdm/gdm-de.omf
/usr/share/omf/gdm/gdm-C.omf
/usr/share/omf/gdm/gdm-uk.omf
/usr/share/omf/gdm/gdm-ko.omf
/usr/share/omf/gdm/gdm-es.omf
/usr/share/omf/gdm/gdm-sv.omf
/usr/share/omf/gdm/gdm-oc.omf
/usr/share/omf/gdm/gdm-en_GB.omf
/usr/share/omf/gdm/gdm-fr.omf
/usr/share/doc
/usr/share/doc/gdm
/usr/share/doc/gdm/README.Debian
/usr/share/doc/gdm/TODO
/usr/share/doc/gdm/changelog.Debian.gz
/usr/share/doc/gdm/README.gz
/usr/share/doc/gdm/AUTHORS.gz
/usr/share/doc/gdm/changelog.gz
/usr/share/doc/gdm/copyright
/usr/share/doc/gdm/NEWS.gz
/usr/share/doc/gdm/README.install.gz
/usr/share/gnome
/usr/share/gnome/help
/usr/share/gnome/help/gdm
/usr/share/gnome/help/gdm/uk
/usr/share/gnome/help/gdm/uk/gdm.xml
/usr/share/gnome/help/gdm/C
/usr/share/gnome/help/gdm/C/legal.xml
/usr/share/gnome/help/gdm/C/gdm.xml
/usr/share/gnome/help/gdm/ko
/usr/share/gnome/help/gdm/ko/gdm.xml
/usr/share/gnome/help/gdm/oc
/usr/share/gnome/help/gdm/oc/gdm.xml
/usr/share/gnome/help/gdm/es
/usr/share/gnome/help/gdm/es/gdm.xml
/usr/share/gnome/help/gdm/en_GB
/usr/share/gnome/help/gdm/en_GB/gdm.xml
/usr/share/gnome/help/gdm/sv
/usr/share/gnome/help/gdm/sv/gdm.xml
/usr/share/gnome/help/gdm/de
/usr/share/gnome/help/gdm/de/gdm.xml
/usr/share/gnome/help/gdm/fr
/usr/share/gnome/help/gdm/fr/gdm.xml
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/24x24
/usr/share/icons/hicolor/24x24/apps
/usr/share/icons/hicolor/24x24/apps/gdmflexiserver.png
/usr/share/icons/hicolor/24x24/apps/gdmsetup.png
/usr/share/icons/hicolor/scalable
/usr/share/icons/hicolor/scalable/apps
/usr/share/icons/hicolor/scalable/apps/gdmflexiserver.svg
/usr/share/icons/hicolor/scalable/apps/gdmsetup.svg
/usr/share/icons/hicolor/32x32
/usr/share/icons/hicolor/32x32/apps
/usr/share/icons/hicolor/32x32/apps/gdmflexiserver.png
/usr/share/icons/hicolor/32x32/apps/gdmsetup.png
/usr/share/icons/hicolor/32x32/apps/gdm-xnest.png
/usr/share/icons/hicolor/16x16
/usr/share/icons/hicolor/16x16/apps
/usr/share/icons/hicolor/16x16/apps/gdmflexiserver.png
/usr/share/icons/hicolor/16x16/apps/gdmsetup.png
/usr/share/icons/hicolor/16x16/apps/gdm-xnest.png
/usr/share/icons/hicolor/48x48
/usr/share/icons/hicolor/48x48/apps
/usr/share/icons/hicolor/48x48/apps/gdm-login-photo.png
/usr/share/icons/hicolor/48x48/apps/gdmflexiserver.png
/usr/share/icons/hicolor/48x48/apps/gdmsetup.png
/usr/share/icons/hicolor/48x48/apps/gdm-xnest.png
/usr/share/icons/hicolor/48x48/apps/gdm.png
/usr/share/icons/hicolor/22x22
/usr/share/icons/hicolor/22x22/apps
/usr/share/icons/hicolor/22x22/apps/gdmflexiserver.png
/usr/share/icons/hicolor/22x22/apps/gdmsetup.png
/usr/share/xsessions
/usr/share/xsessions/ssh.desktop
/usr/share/hosts
/usr/share/faces
/usr/share/menu
/usr/share/menu/gdm
/usr/share/applications
/usr/share/applications/gdmphotosetup.desktop
/usr/share/applications/gdmflexiserver.desktop
/usr/share/applications/gdmsetup.desktop
/usr/share/applications/gdmflexiserver-xnest.desktop
/etc
/etc/gdm
/etc/gdm/gdm.conf-custom
/etc/gdm/Init
/etc/gdm/Init/Default
/etc/gdm/XKeepsCrashing
/etc/gdm/PostLogin
/etc/gdm/PostLogin/Default.sample
/etc/gdm/gdm.conf
/etc/gdm/PreSession
/etc/gdm/PreSession/Default
/etc/gdm/Xsession
/etc/gdm/gdmprefetchlist
/etc/gdm/PostSession
/etc/gdm/PostSession/Default
/etc/gdm/locale.conf
/etc/gdm/modules
/etc/gdm/modules/factory-AccessKeyMouseEvents
/etc/gdm/modules/AccessDwellMouseEvents
/etc/gdm/modules/factory-AccessDwellMouseEvents
/etc/gdm/modules/AccessKeyMouseEvents
/etc/init.d
/etc/init.d/gdm
/etc/pam.d
/etc/pam.d/gdm
/etc/pam.d/gdm-autologin
/etc/dm
/etc/X11
/usr/bin/gdmXnest
