# filesystem/var/lib/dpkg/info/hal.list
/.
/sbin
/var
/var/run
/var/run/hal
/var/lib
/var/lib/hal
/var/cache
/var/cache/hald
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-hal.rules
/etc
/etc/hal
/etc/hal/fdi
/etc/hal/fdi/policy
/etc/hal/fdi/policy/preferences.fdi
/etc/hal/fdi/preprobe
/etc/hal/fdi/information
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/hal.conf
/etc/default
/etc/default/hal
/etc/init.d
/etc/init.d/hal
/usr
/usr/share
/usr/share/hal
/usr/share/hal/fdi
/usr/share/hal/fdi/fdi.dtd
/usr/share/hal/fdi/policy
/usr/share/hal/fdi/policy/20thirdparty
/usr/share/hal/fdi/policy/10osvendor
/usr/share/hal/fdi/policy/10osvendor/10-usbcsr-mice.fdi
/usr/share/hal/fdi/policy/10osvendor/10-input-policy.fdi
/usr/share/hal/fdi/policy/10osvendor/15-storage-luks.fdi
/usr/share/hal/fdi/policy/10osvendor/10-x11-keymap.fdi
/usr/share/hal/fdi/policy/10osvendor/10-cpufreq.fdi
/usr/share/hal/fdi/policy/10osvendor/10-dockstation.fdi
/usr/share/hal/fdi/policy/10osvendor/10-rfkill-switch.fdi
/usr/share/hal/fdi/policy/10osvendor/10-sony-laptop-brightness.fdi
/usr/share/hal/fdi/policy/10osvendor/30-wol.fdi
/usr/share/hal/fdi/policy/10osvendor/10-power-mgmt-policy.fdi
/usr/share/hal/fdi/policy/10osvendor/10-laptop-panel-mgmt-policy.fdi
/usr/share/hal/fdi/policy/10osvendor/20-storage-methods.fdi
/usr/share/hal/fdi/policy/10osvendor/10-tabletPCs.fdi
/usr/share/hal/fdi/policy/10osvendor/10-leds.fdi
/usr/share/hal/fdi/policy/10osvendor/10-keymap.fdi
/usr/share/hal/fdi/policy/10osvendor/debian-storage-policy-ignore-fixed-crypto-drives.fdi
/usr/share/hal/fdi/policy/10osvendor/10-x11-input.fdi
/usr/share/hal/fdi/policy/10osvendor/10-omap-backlight.fdi
/usr/share/hal/fdi/policy/10osvendor/20-acl-management.fdi
/usr/share/hal/fdi/policy/10osvendor/10-smartcardreaders.fdi
/usr/share/hal/fdi/preprobe
/usr/share/hal/fdi/preprobe/20thirdparty
/usr/share/hal/fdi/preprobe/10osvendor
/usr/share/hal/fdi/information
/usr/share/hal/fdi/information/20thirdparty
/usr/share/hal/fdi/information/10freedesktop
/usr/share/PolicyKit
/usr/share/PolicyKit/policy
/usr/share/PolicyKit/policy/org.freedesktop.hal.wol.policy
/usr/share/PolicyKit/policy/org.freedesktop.hal.policy
/usr/share/PolicyKit/policy/org.freedesktop.hal.leds.policy
/usr/share/PolicyKit/policy/org.freedesktop.hal.storage.policy
/usr/share/PolicyKit/policy/org.freedesktop.hal.killswitch.policy
/usr/share/PolicyKit/policy/org.freedesktop.hal.dockstation.policy
/usr/share/PolicyKit/policy/org.freedesktop.hal.power-management.policy
/usr/share/PolicyKit/policy/org.freedesktop.hal.device-access.policy
/usr/share/doc
/usr/share/doc/hal
/usr/share/doc/hal/examples
/usr/share/doc/hal/examples/no-cd-media-check.fdi
/usr/share/doc/hal/copyright
/usr/share/doc/hal/README.Debian
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/hal-find-by-capability.1.gz
/usr/share/man/man1/hal-is-caller-privileged.1.gz
/usr/share/man/man1/hal-set-property.1.gz
/usr/share/man/man1/hal-disable-polling.1.gz
/usr/share/man/man1/hal-get-property.1.gz
/usr/share/man/man1/hal-is-caller-locked-out.1.gz
/usr/share/man/man1/hal-lock.1.gz
/usr/share/man/man1/hal-find-by-property.1.gz
/usr/share/man/man1/lshal.1.gz
/usr/share/man/man8
/usr/share/man/man8/hald.8.gz
/usr/sbin
/usr/sbin/hald
/usr/sbin/umount.hal
/usr/bin
/usr/bin/hal-device
/usr/bin/hal-find-by-property
/usr/bin/hal-get-property
/usr/bin/hal-setup-keymap
/usr/bin/hal-find-by-capability
/usr/bin/hal-disable-polling
/usr/bin/hal-set-property
/usr/bin/hal-lock
/usr/bin/hal-is-caller-privileged
/usr/bin/lshal
/usr/bin/hal-is-caller-locked-out
/usr/lib
/usr/lib/hal
/usr/lib/hal/hal-storage-mount
/usr/lib/hal/hald-probe-input
/usr/lib/hal/hald-probe-hiddev
/usr/lib/hal/hald-addon-ipw-killswitch
/usr/lib/hal/hald-addon-input
/usr/lib/hal/debian-setup-keyboard
/usr/lib/hal/hal-system-sonypic
/usr/lib/hal/hald-addon-cpufreq
/usr/lib/hal/hal-system-power-pmu
/usr/lib/hal/hald-addon-rfkill-killswitch
/usr/lib/hal/hal-storage-unmount
/usr/lib/hal/hald-probe-serial
/usr/lib/hal/hald-generate-fdi-cache
/usr/lib/hal/hald-addon-generic-backlight
/usr/lib/hal/hald-runner
/usr/lib/hal/hald-probe-ieee1394-unit
/usr/lib/hal/hald-probe-net-bluetooth
/usr/lib/hal/hal-system-setserial
/usr/lib/hal/hald-addon-omap-backlight
/usr/lib/hal/hal-storage-cleanup-mountpoint
/usr/lib/hal/hal-system-power-pm-is-supported
/usr/lib/hal/hald-probe-smbios
/usr/lib/hal/scripts
/usr/lib/hal/scripts/hal-system-lcd-get-brightness
/usr/lib/hal/scripts/hal-system-killswitch-get-power
/usr/lib/hal/scripts/hal-functions
/usr/lib/hal/scripts/hal-system-power-set-power-save
/usr/lib/hal/scripts/hal-system-power-reboot
/usr/lib/hal/scripts/hal-system-wol-enable
/usr/lib/hal/scripts/hal-luks-teardown
/usr/lib/hal/scripts/linux
/usr/lib/hal/scripts/linux/hal-system-power-set-power-save-linux
/usr/lib/hal/scripts/linux/hal-system-power-suspend-linux
/usr/lib/hal/scripts/linux/hal-system-killswitch-get-power-linux
/usr/lib/hal/scripts/linux/hal-luks-teardown-linux
/usr/lib/hal/scripts/linux/hal-system-lcd-get-brightness-linux
/usr/lib/hal/scripts/linux/hal-luks-setup-linux
/usr/lib/hal/scripts/linux/hal-system-lcd-set-brightness-linux
/usr/lib/hal/scripts/linux/hal-system-power-shutdown-linux
/usr/lib/hal/scripts/linux/hal-system-power-hibernate-linux
/usr/lib/hal/scripts/linux/hal-system-killswitch-set-power-linux
/usr/lib/hal/scripts/linux/hal-system-power-suspend-hybrid-linux
/usr/lib/hal/scripts/linux/hal-luks-remove-linux
/usr/lib/hal/scripts/linux/hal-system-wol-linux
/usr/lib/hal/scripts/linux/hal-system-power-reboot-linux
/usr/lib/hal/scripts/linux/hal-dockstation-undock-linux
/usr/lib/hal/scripts/hal-system-power-suspend
/usr/lib/hal/scripts/hal-system-power-hibernate
/usr/lib/hal/scripts/hal-system-power-suspend-hybrid
/usr/lib/hal/scripts/hal-system-wol-supported
/usr/lib/hal/scripts/hal-system-wol-enabled
/usr/lib/hal/scripts/hal-system-power-shutdown
/usr/lib/hal/scripts/hal-luks-setup
/usr/lib/hal/scripts/hal-luks-remove
/usr/lib/hal/scripts/hal-dockstation-undock
/usr/lib/hal/scripts/hal-system-killswitch-set-power
/usr/lib/hal/scripts/hal-system-lcd-set-brightness
/usr/lib/hal/hald-addon-storage
/usr/lib/hal/hald-addon-acpi
/usr/lib/hal/hal-acl-tool
/usr/lib/hal/hald-addon-leds
/usr/lib/hal/hald-probe-printer
/usr/lib/hal/hal-storage-closetray
/usr/lib/hal/hal-storage-cleanup-all-mountpoints
/usr/lib/hal/hald-addon-pmu
/usr/lib/hal/hald-probe-storage
/usr/lib/hal/hal-storage-eject
/usr/lib/hal/hald-probe-pc-floppy
/usr/lib/hal/hald-probe-video4linux
/usr/lib/hal/hald-probe-volume
/usr/lib/hal/hald-addon-usb-csr
/usr/lib/hal/hald-addon-hid-ups
/sbin/umount.hal
/usr/share/doc/hal/changelog.Debian.gz
/usr/share/doc/hal/README
/usr/share/doc/hal/NEWS.gz
/usr/share/doc/hal/AUTHORS
/usr/lib/hal/scripts/linux/hal-system-wol-supported-linux
/usr/lib/hal/scripts/linux/hal-system-wol-enable-linux
/usr/lib/hal/scripts/linux/hal-system-wol-enabled-linux
