# filesystem/var/lib/dpkg/info/netbase.conffiles
/etc/protocols
/etc/rpc
/etc/services
/etc/init.d/networking
