# filesystem/var/lib/dpkg/info/policykit.conffiles
/etc/init.d/policykit
/etc/pam.d/polkit
/etc/PolicyKit/PolicyKit.conf
/etc/dbus-1/system.d/org.freedesktop.PolicyKit.conf
