# filesystem/var/lib/dpkg/info/policykit.list
/.
/usr
/usr/bin
/usr/bin/polkit-auth
/usr/bin/polkit-config-file-validate
/usr/bin/polkit-policy-file-validate
/usr/bin/polkit-action
/usr/lib
/usr/lib/policykit
/usr/lib/policykit/polkit-read-auth-helper
/usr/lib/policykit/polkit-set-default-helper
/usr/lib/policykit/polkit-resolve-exe-helper
/usr/lib/policykit/polkit-grant-helper
/usr/lib/policykit/polkit-grant-helper-pam
/usr/lib/policykit/polkit-revoke-helper
/usr/lib/policykit/polkit-explicit-grant-helper
/usr/lib/policykit/polkitd
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/PolicyKit.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/PolicyKit.8.gz
/usr/share/man/man1
/usr/share/man/man1/polkit-policy-file-validate.1.gz
/usr/share/man/man1/polkit-auth.1.gz
/usr/share/man/man1/polkit-config-file-validate.1.gz
/usr/share/man/man1/polkit-action.1.gz
/usr/share/PolicyKit
/usr/share/PolicyKit/config.dtd
/usr/share/PolicyKit/policy
/usr/share/PolicyKit/policy/org.freedesktop.policykit.policy
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.PolicyKit.service
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.PolicyKit.AuthenticationAgent.xml
/usr/share/doc
/usr/share/doc/policykit
/usr/share/doc/policykit/copyright
/var
/var/lib
/var/lib/PolicyKit
/var/lib/misc
/var/lib/misc/PolicyKit.reload
/var/lib/PolicyKit-public
/etc
/etc/init.d
/etc/init.d/policykit
/etc/pam.d
/etc/pam.d/polkit
/etc/PolicyKit
/etc/PolicyKit/PolicyKit.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.PolicyKit.conf
/usr/share/doc/policykit/README
/usr/share/doc/policykit/NEWS.gz
/usr/share/doc/policykit/changelog.Debian.gz
/usr/share/doc/policykit/AUTHORS
