# filesystem/var/lib/dpkg/info/procps.conffiles
/etc/sysctl.d/10-network-security.conf
/etc/sysctl.d/README
/etc/sysctl.d/10-console-messages.conf
/etc/init.d/procps
/etc/sysctl.conf
