# filesystem/var/lib/dpkg/info/udev.list
/.
/sbin
/sbin/udevd
/sbin/udevadm
/lib
/lib/firmware
/lib/udev
/lib/udev/edd_id
/lib/udev/create_floppy_devices
/lib/udev/devices
/lib/udev/devices/pts
/lib/udev/devices/shm
/lib/udev/devices/net
/lib/udev/rules.d
/lib/udev/rules.d/60-persistent-storage-tape.rules
/lib/udev/rules.d/95-udev-late.rules
/lib/udev/rules.d/60-cdrom_id.rules
/lib/udev/rules.d/50-firmware.rules
/lib/udev/rules.d/60-persistent-v4l.rules
/lib/udev/rules.d/61-persistent-storage-edd.rules
/lib/udev/rules.d/40-pilot-links.rules
/lib/udev/rules.d/40-ia64.rules
/lib/udev/rules.d/80-drivers.rules
/lib/udev/rules.d/60-persistent-storage.rules
/lib/udev/rules.d/40-isdn.rules
/lib/udev/rules.d/50-udev-default.rules
/lib/udev/rules.d/79-fstab_import.rules
/lib/udev/rules.d/README
/lib/udev/rules.d/60-persistent-serial.rules
/lib/udev/rules.d/40-zaptel.rules
/lib/udev/rules.d/40-alsa.rules
/lib/udev/rules.d/40-infiniband.rules
/lib/udev/rules.d/75-cd-aliases-generator.rules
/lib/udev/rules.d/60-persistent-input.rules
/lib/udev/rules.d/40-ppc.rules
/lib/udev/rules.d/64-device-mapper.rules
/lib/udev/rules.d/75-persistent-net-generator.rules
/lib/udev/vol_id
/lib/udev/write_cd_rules
/lib/udev/write_net_rules
/lib/udev/ata_id
/lib/udev/scsi_id
/lib/udev/usb_id
/lib/udev/firmware.sh
/lib/udev/fstab_import
/lib/udev/path_id
/lib/udev/cdrom_id
/lib/udev/rule_generator.functions
/lib/udev/collect
/etc
/etc/scsi_id.config
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/README
/etc/udev/udev.conf
/etc/init.d
/etc/init.d/udev-finish
/etc/init.d/udev
/usr
/usr/share
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/udev
/usr/share/initramfs-tools/scripts/nfs-top
/usr/share/initramfs-tools/scripts/nfs-top/udev
/usr/share/initramfs-tools/scripts/local-top
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/udev
/usr/share/initramfs-tools/scripts/casper-premount
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/udev
/usr/share/doc
/usr/share/doc/udev
/usr/share/doc/udev/examples
/usr/share/doc/udev/examples/scsi_id.config
/usr/share/doc/udev/copyright
/usr/share/doc/udev/changelog.gz
/usr/share/doc/udev/changelog.Debian.gz
/usr/share/doc/udev/writing_udev_rules
/usr/share/doc/udev/writing_udev_rules/index.html
/usr/share/doc/udev/README
/usr/share/doc/udev/NEWS.gz
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/udev.7.gz
/usr/share/man/man8
/usr/share/man/man8/udevadm.8.gz
/usr/share/man/man8/udevd.8.gz
/usr/share/man/man8/edd_id.8.gz
/usr/share/man/man8/scsi_id.8.gz
/usr/share/man/man8/path_id.8.gz
/usr/share/man/man8/usb_id.8.gz
/usr/share/man/man8/ata_id.8.gz
/usr/share/man/man8/cdrom_id.8.gz
/usr/share/man/man8/vol_id.8.gz
/usr/share/man/man8/create_floppy_devices.8.gz
/usr/lib
/usr/lib/udev
/sbin/vol_id
